# supervisor.conf
# VinApp - Supervisor config for Laravel queue workers
# Forge lager dette automatisk via "Daemons"-seksjonen i dashboardet
# Legg til via: Forge → Server → Daemons → New Daemon
#
# Command:  php /home/forge/vinapp.no/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
# User:     forge
# Processes: 2

[program:vinapp-worker]
process_name=%(program_name)s_%(process_unm)02d
command=php /home/forge/vinapp.no/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=2
redirect_stderr=true
stdout_logfile=/home/forge/.forge/vinapp-worker.log
stopwaitsecs=3600
